# recipes-debian/busybox/files/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
